# scripts/systemd/director-beat.service
[Unit]
Description=Directely Celery beat (scheduled tasks)
After=director-infra.service network-online.target
Requires=director-infra.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/director
ExecStart=/root/director/scripts/systemd/run-beat.sh
Restart=on-failure
RestartSec=5
TimeoutStartSec=120

[Install]
WantedBy=multi-user.target
